# packages/PyBarCodeScan/PyBarCodeScan-0.0.10.tar.gz/PyBarCodeScan-0.0.10/barcodescan/barcodescan_sse.service
[Unit]
Description=barcodescan_sse

[Service]
Type=simple
ExecStart=/venv/bin/barcodescan_sse

[Install]
WantedBy=multi-user.target