# packages/PyBarCodeScan/PyBarCodeScan-0.1.48.tar.gz/PyBarCodeScan-0.1.48/barcodescan/barcodescanner_reader.service
[Unit]
Description=Reads barcodes from hardware and send them to barcodescanner_server

[Service]
Type=simple
ExecStart=/venv/bin/barcodescanner_reader

[Install]
WantedBy=multi-user.target