# packages/tendril-server-fs/tendril-server-fs-0.2.9.tar.gz/tendril-server-fs-0.2.9/deploy/tendril_fs_server.service

[Unit]
Description=Tendril FS Server

[Service]
ExecStart=/usr/bin/twistd --nodaemon --pidfile= fs_server.tac

WorkingDirectory=/home/tendril/tools/fs_server/deploy

User=tendril
Group=tendril

Restart=always

[Install]
WantedBy=multi-user.target

