# packages/PyBarCodeScan/PyBarCodeScan-0.0.8.tar.gz/PyBarCodeScan-0.0.8/barcodescan/barcodescan_sse.service
[Unit]
Description=barcodescan_sse

[Service]
Type=simple
ExecStart=barcodescan_sse

[Install]
WantedBy=multi-user.target