# packages/PyBarCodeScan/PyBarCodeScan-0.0.12.tar.gz/PyBarCodeScan-0.0.12/barcodescan/read_from_scanner.service
[Unit]
Description=read_from_scanner

[Service]
Type=simple
ExecStart=/venv/bin/read_from_scanner

[Install]
WantedBy=multi-user.target